# primer.service
[Unit]
Description=Primer AI Kit Service
Wants=network-online.target
After=network-online.target ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/Primer-Software
Environment="PATH=/home/ubuntu/projects/Primer-Software/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/ubuntu"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/projects/Primer-Software/venv/bin/python /home/ubuntu/projects/Primer-Software/src/primer.py
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target